# deploy/asterisk-mqtt.service
[Unit]
Description=Asterisk-MQTT Bridge
Documentation=https://github.com/sweeney/asterisk-mqtt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=asterisk-mqtt
Group=asterisk-mqtt
ExecStart=/usr/local/bin/asterisk-mqtt -config /etc/asterisk-mqtt/asterisk-mqtt.yaml
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/etc/asterisk-mqtt

[Install]
WantedBy=multi-user.target
